# scripts/ZTBWIFI/files/etc/uci-defaults/10_disable_services
/etc/init.d/homeproxy disable
/etc/init.d/p910nd disable
/etc/init.d/parentcontrol disable
/etc/init.d/timewol disable

exit 0