# files/etc/uci-defaults/10_disable_services
/etc/init.d/parentcontrol disable
/etc/init.d/timewol disable
/etc/init.d/weburl disable

exit 0